# com.jack.ktbfx.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jack.ktbfx.server</string>

  <!-- Keeps the local collector endpoint (server.js, POST /mc) running so the
       Mastercard capture browser extension always has something to post to.
       Also serves the dashboard at http://localhost:8777 . -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>export PATH="/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin"; exec node server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jack/fx-dashboard</string>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key>
  <string>/Users/jack/fx-dashboard/logs/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jack/fx-dashboard/logs/server.err.log</string>
</dict>
</plist>
